# package/fcm/fwxd/files/uci-defaults/100_fwx
#!/bin/sh
uci set fwx.global.theme_mode='1'
uci commit fwx
